# .config/systemd/user/backup.service
[Unit]
Description = Run backup script

[Service]
Type=oneshot
ExecStart=/home/polhuang/.dotfiles/.bin/backup.sh 

[Install]
WantedBy=default.target